# hostlogger.service.in
[Unit]
Description=Host logging
Wants=obmc-console@@HOST_TTY@.service
After=obmc-console@@HOST_TTY@.service

[Service]
ExecStart=hostlogger
Restart=always
Type=dbus
BusName=xyz.openbmc_project.HostLogger

[Install]
WantedBy=@SYSTEMD_TARGET@
